# systemd/endoscope-viewer.service
[Unit]
Description=USB Endoscope → TFT viewer (simple)
After=multi-user.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__WORKDIR__/scripts/run_service.sh --rotate 90
Restart=always
RestartSec=1

[Install]
WantedBy=multi-user.target
